# webcd.service
[Unit]
Description=WebCD - Web-based CD Player
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Code/webcd
Environment="PATH=/home/%i/Code/webcd/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/%i/Code/webcd/venv/bin/python /home/%i/Code/webcd/app.py
Restart=on-failure
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/Code/webcd
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Device access
DeviceAllow=/dev/cdrom rw
DeviceAllow=/dev/sr0 rw
PrivateDevices=false

[Install]
WantedBy=multi-user.target