# systemd/webcd.service
[Unit]
Description=WebCD - Web-based CD Player
Documentation=https://github.com/GlassOnTin/webcd
After=network.target sound.target

[Service]
Type=simple
User=webcd
Group=cdrom
WorkingDirectory=/opt/webcd
Environment="PATH=/opt/webcd/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/webcd"
ExecStart=/opt/webcd/venv/bin/python /opt/webcd/app.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp
ReadOnlyPaths=/opt/webcd
# Allow access to CD devices
DeviceAllow=/dev/cdrom rw
DeviceAllow=/dev/dvd rw
DeviceAllow=/dev/sr0 rw
DeviceAllow=/dev/sr1 rw
DeviceAllow=/dev/sr2 rw
DevicePolicy=closed
SupplementaryGroups=cdrom audio

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target