# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create webcd system user if it doesn't exist
        if ! getent passwd webcd > /dev/null; then
            adduser --system --group --home /nonexistent --no-create-home \
                --quiet --disabled-login --disabled-password \
                --gecos "WebCD CD Player" webcd
        fi
        
        # Add webcd user to cdrom and audio groups
        if getent group cdrom > /dev/null; then
            adduser webcd cdrom || true
        fi
        if getent group audio > /dev/null; then
            adduser webcd audio || true
        fi
        
        # Set permissions on application directory
        chown -R webcd:webcd /usr/share/webcd
        chmod 755 /usr/share/webcd
        
        # Reload systemd and enable service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
            # Enable the service to start on boot
            systemctl enable webcd.service || true
            # Start the service
            systemctl start webcd.service || true
        fi
        
        echo "WebCD installed successfully!"
        echo "Service has been started and enabled."
        echo "Access WebCD at: http://localhost:5000"
        echo ""
        echo "Useful commands:"
        echo "  systemctl status webcd    # Check service status"
        echo "  systemctl restart webcd   # Restart service"
        echo "  journalctl -u webcd -f    # View logs"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0